# tools/clear_img_file.cron
* * * * * find  /opt/dianshu/www/weixin/static/douban_image/ atime +1 | xargs rm -f 
